# lgfc-traffic-simulator.service
[Unit]
Description=Lou Gehrig Fan Club Traffic Simulator
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/workspace
ExecStart=/usr/bin/node /workspace/lgfc_enhanced_jwt_traffic_simulator.cjs --interval=300000 --users=15
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lgfc-traffic-simulator

# Restart policy
StartLimitInterval=60
StartLimitBurst=3

# Environment
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target